# etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
Type=simple
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \u' --noclear --autologin rakib %I xterm-256color
