# etc/systemd/system/bluetooth.service.d/bluetooth.conf
# --experimental enables airpods battery percentage.
[Service]
ExecStart=
ExecStart=/usr/lib/bluetooth/bluetoothd --experimental --kernel
